# SRC/jconf/toad
#!/bin/sh
. /etc/rc.subr
name="toad"
rcvar=toad_enable
load_rc_config $name
/usr/local/bin/hypnotoad /toad/TOAD
